# project-template/config/prod/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb

nginx -g "daemon off;"
